# service/seafile/run
#!/bin/bash -e

set -eu

exec 2>&1
exec ./seafile-start.sh
