# service/seahub/run
#!/bin/bash

set -eu

exec 2>&1
exec ./seahub-start.sh
